# deployment/conductor-supervisor.conf
[program:conductor]
command=/home/josh/.nix-profile/bin/holochain --config /home/josh/Junto-Holochain-Alpha-API/config.toml
user=josh