# deployment/conductor-supervisor.conf
[program:conductor]
command=holochain --config /home/josh/Junto-Holochain-Alpha-API/conductor_generator/config.toml
file=/tmp/supervisor.sock
chmod=0766