# deployment/web-app-supervisor.conf
[program:webapp]
command=bash -c 'cd ~/Junto-Holochain-Alpha-API && /home/josh/.cargo/bin/cargo run --release'
user=josh